# debian/cron.d
# Delete all old containers ...
0 1 * * * joulupukki docker rm $(docker ps -a |grep joulupukki | grep -vE "minutes|hours" | awk '{print $1}')
